# pir-watcher/pir-watcher.service
[Unit]
Description=PIR motion watcher -> /tmp slot + S3 PUT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/pir-monitor/config.env
# NTP同期完了（NTPSynchronized=yes）まで待機（最大60秒ならsleep回数を調整）
ExecStartPre=/bin/bash -lc 'for i in {1..30}; do \
  timedatectl show -p NTPSynchronized --value | grep -q yes && exit 0; \
  sleep 2; \
done; exit 0'
ExecStart=/usr/bin/python3 /usr/local/bin/pir-watcher.py
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target