# sync_device_config/sync_device_config.timer
[Unit]
Description=Run sync device config

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target