# sync_device_config/sync_device_config.service
[Unit]
Description=Sync device config

[Service]
Type=oneshot
EnvironmentFile=/etc/pir-monitor/config.env
ExecStart=/usr/bin/python3 /usr/local/bin/sync_device_config.py
User=root